# cc-sage-proxy/extra/etc/systemd/system/setup-firewall.service
[Unit]
Description=Setup firewall and routing

[Service]
ExecStart=/root/setup-firewall

[Install]
WantedBy=multi-user.target
